# decompressed/gui_file/etc/crontabs/root_gui
0 5 * * * /usr/share/transformer/scripts/checkver.sh Schedule >/dev/null 2>&1