# decompressed/upgrade-pack-specificDGA/usr/lib/opkg/info/nginx.list
/etc/nginx/fastcgi_params
/etc/nginx/naxsi_core.rules
/etc/nginx/koi-win
/usr/sbin/nginx
/etc/nginx/win-utf
/etc/nginx/koi-utf
/etc/init.d/nginx
/etc/nginx/mime.types
/lib/upgrade/keep.d/nginx
