# etc/crontabs/root_gui
0 */2 * * * /usr/share/transformer/scripts/checkver.sh